# gimbal-tracking.service
[Unit]
Description=Gimbal Tracking Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linaro
Group=linaro
SupplementaryGroups=dialout
WorkingDirectory=/home/linaro/gimbal
EnvironmentFile=-/etc/default/gimbal-tracking
ExecStart=/home/linaro/gimbal/run_main_tracking.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
